# deploy/systemd/aped.service
# The root executor (PLAN-18 D1 / Appendix A): root without power. Empty
# capability set (non-device tier — containerd/shim/QEMU hold every VM
# capability); AF_UNIX only, no network; @system-service syscall filter (already
# allows the KVM/VFIO ioctl + chown). Measure with:
#   systemd-analyze security aped.service   (predicted OK band ~3.0–3.8)
#
# NOTE: systemd has NO inline comments — a trailing "# ..." becomes part of the
# value. Every comment here is on its own line for that reason.
[Unit]
Description=aped — rootful Kata-QEMU VM-management executor (privileged, narrow)
After=containerd.service network-pre.target
Requires=aped-priv.socket
Wants=containerd.service

[Service]
# Type=notify: aped signals READY=1 once the priv socket is bound/adopted and
# pings the watchdog on a timer (sd_notify(3), via $NOTIFY_SOCKET — AF_UNIX, so
# RestrictAddressFamilies=AF_UNIX below already permits it). The priv socket is
# inherited from aped-priv.socket (LISTEN_FDS, auto-detected).
Type=notify
ExecStart=/usr/local/bin/aped run --policy /etc/aped/policy.yaml --state-dir /var/lib/aped --audit-log /var/log/aped/audit.jsonl --node %H --allow-user aped
Restart=on-failure
WatchdogSec=30s

# Identity: root, for the containerd socket + (device tier) sysfs vfio-pci bind.
# Group=ape is for the priv.sock group; NOT the containerd group.
User=root
Group=ape
NoNewPrivileges=yes

# Capabilities: empty (non-device). Device tier adds only CAP_CHOWN (see PLAN-18 D4).
CapabilityBoundingSet=
AmbientCapabilities=

# Filesystem. StateDirectoryMode=0770 keeps /var/lib/aped root:ape group-writable
# so the de-privileged front (User=aped, group ape) shares this tree.
StateDirectory=aped
StateDirectoryMode=0770
LogsDirectory=aped
RuntimeDirectory=aped
RuntimeDirectoryMode=0710
# PRESERVE /run/aped across restarts. aped-priv.socket binds /run/aped/priv.sock
# there; without this the default RuntimeDirectoryPreserve=no makes systemd remove
# /run/aped when aped.service STOPS, taking the priv socket file with it — so a
# `restart aped-priv.socket` then `restart aped.service` desyncs the socket
# (front: "dial /run/aped/priv.sock: no such file"). tmpfiles.d owns the dir's
# existence; this just stops the service teardown from deleting it.
RuntimeDirectoryPreserve=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/aped /run/aped /var/log/aped
ProtectHome=yes
PrivateTmp=yes
UMask=0077
ProtectProc=invisible
ProcSubset=pid
# ProtectKernelTunables — DEVICE tier omits this (vfio-pci bind writes /sys).
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
PrivateDevices=yes

# Namespaces: aped creates none (containerd/shim do). Deny all.
RestrictNamespaces=yes

# Network: AF_UNIX only — unreachable over IP.
RestrictAddressFamilies=AF_UNIX
IPAddressDeny=any

# Syscalls: @system-service already allows ioctl (KVM/VFIO) + chown + process.
# The second line is a single denylist (one ~ prefixes the whole space-separated
# group list — writing ~ before each group makes systemd read "~@swap" as an
# unknown syscall NAME and silently drop the filter).
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM
SystemCallFilter=@system-service
SystemCallFilter=~@mount @swap @reboot @module @raw-io @cpu-emulation @obsolete @debug @privileged

RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
KeyringMode=private

[Install]
WantedBy=multi-user.target
