# deploy/systemd/aped.service
# The root executor (PLAN-18 D1 / Appendix A): root without power. Empty
# capability set (non-device tier — containerd/shim/QEMU hold every VM
# capability); AF_UNIX only, no network; @system-service syscall filter (already
# allows the KVM/VFIO ioctl + chown). Measure with:
#   systemd-analyze security aped.service   (predicted OK band ~3.0–3.8)
[Unit]
Description=aped — rootful Kata-QEMU VM-management executor (privileged, narrow)
After=containerd.service network-pre.target
Requires=aped-priv.socket
Wants=containerd.service

[Service]
# Type=exec (not notify): aped does not yet implement sd_notify. The priv socket
# is inherited from aped-priv.socket (LISTEN_FDS, auto-detected).
Type=exec
ExecStart=/usr/bin/aped run --policy /etc/aped/policy.yaml --state-dir /var/lib/aped --audit-log /var/log/aped/audit.jsonl --node %H --allow-user aped
Restart=on-failure

# Identity: root, for the containerd socket + (device tier) sysfs vfio-pci bind.
User=root
Group=ape                     # for the priv.sock group; NOT the containerd group
NoNewPrivileges=yes

# Capabilities: empty (non-device). Device tier adds only CAP_CHOWN (see PLAN-18 D4).
CapabilityBoundingSet=
AmbientCapabilities=

# Filesystem
StateDirectory=aped
LogsDirectory=aped
RuntimeDirectory=aped
ProtectSystem=strict
ReadWritePaths=/var/lib/aped /run/aped /var/log/aped
ProtectHome=yes
PrivateTmp=yes
UMask=0077
ProtectProc=invisible
ProcSubset=pid
ProtectKernelTunables=yes      # DEVICE tier omits this (vfio-pci bind writes /sys)
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
PrivateDevices=yes

# Namespaces: aped creates none (containerd/shim do). Deny all.
RestrictNamespaces=yes

# Network: AF_UNIX only — unreachable over IP.
RestrictAddressFamilies=AF_UNIX
IPAddressDeny=any

# Syscalls: @system-service already allows ioctl (KVM/VFIO) + chown + process.
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM
SystemCallFilter=@system-service
SystemCallFilter=~@mount ~@swap ~@reboot ~@module ~@raw-io ~@cpu-emulation ~@obsolete ~@debug ~@privileged

RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
KeyringMode=private

[Install]
WantedBy=multi-user.target
